# src/crontab
* * * * * root . /root/cron.env; /usr/local/bin/node -e "console.log(process.env)" >> /var/log/env.cron.log 2>&1 

